# root/etc/services.d/overseerr/run
#!/bin/bash
set -euo pipefail

export CONFIG_DIRECTORY=/config/overseerr

(
	cd /overseerr
	npm start
)
