# root/etc/services.d/rclone/run
#!/usr/bin/env bash
set -euo pipefail

mkdir -p /config/rclone

if [ ! -f /config/rclone/start.sh ]; then
	cp /templates/rclone/start.sh /config/rclone/start.sh
	echo "#######################################################################"
	echo "Copied rclone start script to /config/rclone/start.sh."
	echo "#######################################################################"
fi

echo "Starting rclone."
rm -rf /media
mkdir -p /media
/config/rclone/start.sh
