# root/etc/services.d/connector/run
#!/bin/bash
set -euo pipefail

export STORAGE_DIR=/config/connector

(
	cd /app/connector
	pnpm start-in-coaxist
)
