# root/etc/services.d/connector/run
#!/bin/bash
set -euo pipefail

export STORAGE_DIR=/config/connector
export OVERSEERR_CONFIG_PATH=/config/overseerr/settings.json

(
	cd /app/connector
	pnpm start-in-coaxist
)
